# files/etc/uci-defaults/79_generated_settings
#!/bin/ash
uci set alfred.alfred=alfred
uci set alfred.alfred.interface='br-guest'
uci set alfred.alfred.mode='slave'
uci set alfred.alfred.batmanif='bat0'
uci set alfred.alfred.start_vis='1'
uci set alfred.alfred.run_facters='0'
uci set batman-adv.bat0=mesh
uci set batman-adv.bat0.multicast_mode='0'
uci set batman-adv.bat0.hop_penalty='15'
uci set batman-adv.bat0.gw_mode='client'
uci set batman-adv.bat0.orig_interval='5000'
uci set dhcp.lan=dhcp
uci set dhcp.lan.interface='lan'
uci set dhcp.lan.dhcpv6='server'
uci set dhcp.lan.ra='server'
uci set dhcp.lan.ra_management='1'
uci set dhcp.lan.ignore='1'
uci set dhcp.wan=dhcp
uci set dhcp.wan.interface='wan'
uci set dhcp.wan.ignore='1'
uci set dhcp.@dnsmasq[0]=dnsmasq
uci set dhcp.@dnsmasq[0].domainneeded='1'
uci set dhcp.@dnsmasq[0].boguspriv='1'
uci set dhcp.@dnsmasq[0].filterwin2k='0'
uci set dhcp.@dnsmasq[0].localise_queries='1'
uci set dhcp.@dnsmasq[0].rebind_protection='1'
uci set dhcp.@dnsmasq[0].rebind_localhost='1'
uci set dhcp.@dnsmasq[0].local='/lan/'
uci set dhcp.@dnsmasq[0].expandhosts='1'
uci set dhcp.@dnsmasq[0].nonegcache='0'
uci set dhcp.@dnsmasq[0].authoritative='1'
uci set dhcp.@dnsmasq[0].readethers='1'
uci set dhcp.@dnsmasq[0].leasefile='/tmp/dhcp.leases'
uci set dhcp.@dnsmasq[0].resolvfile='/tmp/resolv.conf.auto'
uci set dhcp.@dnsmasq[0].localservice='1'
uci set dhcp.@dnsmasq[0].domain='ffcp.lan'
uci set dropbear.@dropbear[0]=dropbear
uci set dropbear.@dropbear[0].PasswordAuth='on'
uci set dropbear.@dropbear[0].Port='22'
uci set fastd.sample_config=fastd
uci set fastd.sample_config.enabled='0'
uci set fastd.sample_config.syslog_level='info'
uci set fastd.sample_config.method='salsa2012+umac'
uci set fastd.sample_config.mode='tap'
uci set fastd.sample_config.interface='tap0'
uci set fastd.sample_config.mtu='1426'
uci set fastd.sample_config.forward='0'
uci set fastd.sample_config.secure_handshakes='1'
uci set fastd.sample_peer=peer
uci set fastd.sample_peer.enabled='0'
uci set fastd.sample_peer.net='sample_config'
uci set fastd.sample_peer.key='0000000000000000000000000000000000000000000000000000000000000000'
uci set fastd.sample_group=peer_group
uci set fastd.sample_group.enabled='0'
uci set fastd.sample_group.net='sample_config'
uci set firewall.@defaults[0]=defaults
uci set firewall.@defaults[0].syn_flood='1'
uci set firewall.@defaults[0].input='ACCEPT'
uci set firewall.@defaults[0].output='ACCEPT'
uci set firewall.@defaults[0].forward='REJECT'
uci set firewall.@rule[2]=rule
uci set firewall.@rule[2].name='Allow-DHCPv6'
uci set firewall.@rule[2].src='wan'
uci set firewall.@rule[2].proto='udp'
uci set firewall.@rule[2].src_ip='fe80::/10'
uci set firewall.@rule[2].src_port='547'
uci set firewall.@rule[2].dest_ip='fe80::/10'
uci set firewall.@rule[2].dest_port='546'
uci set firewall.@rule[2].family='ipv6'
uci set firewall.@rule[2].target='ACCEPT'
uci set firewall.@rule[0]=rule
uci set firewall.@rule[0].name='Allow-DHCP-Renew'
uci set firewall.@rule[0].src='wan'
uci set firewall.@rule[0].proto='udp'
uci set firewall.@rule[0].dest_port='68'
uci set firewall.@rule[0].target='ACCEPT'
uci set firewall.@rule[0].family='ipv4'
uci set firewall.@rule[4]=rule
uci set firewall.@rule[4].name='Allow-ICMPv6-Forward'
uci set firewall.@rule[4].src='wan'
uci set firewall.@rule[4].dest='*'
uci set firewall.@rule[4].proto='icmp'
uci set firewall.@rule[4].icmp_type='echo-request' 'echo-reply' 'destination-unreachable' 'packet-too-big' 'time-exceeded' 'bad-header' 'unknown-header-type'
uci set firewall.@rule[4].limit='1000/sec'
uci set firewall.@rule[4].family='ipv6'
uci set firewall.@rule[3]=rule
uci set firewall.@rule[3].name='Allow-ICMPv6-Input'
uci set firewall.@rule[3].src='wan'
uci set firewall.@rule[3].proto='icmp'
uci set firewall.@rule[3].icmp_type='echo-request' 'echo-reply' 'destination-unreachable' 'packet-too-big' 'time-exceeded' 'bad-header' 'unknown-header-type' 'router-solicitation' 'neighbour-solicitation' 'router-advertisement' 'neighbour-advertisement'
uci set firewall.@rule[3].limit='1000/sec'
uci set firewall.@rule[3].family='ipv6'
uci set firewall.@rule[3].target='ACCEPT'
uci set firewall.@rule[1]=rule
uci set firewall.@rule[1].name='Allow-Ping'
uci set firewall.@rule[1].src='wan'
uci set firewall.@rule[1].proto='icmp'
uci set firewall.@rule[1].icmp_type='echo-request'
uci set firewall.@rule[1].family='ipv4'
uci set firewall.@rule[1].target='ACCEPT'
uci set firewall.@zone[2]=zone
uci set firewall.@zone[2].forward='REJECT'
uci set firewall.@zone[2].name='guest'
uci set firewall.@zone[2].input='REJECT'
uci set firewall.@zone[2].output='REJECT'
uci set firewall.@zone[2].network='guest'
uci set firewall.@zone[0]=zone
uci set firewall.@zone[0].name='lan'
uci set firewall.@zone[0].input='ACCEPT'
uci set firewall.@zone[0].output='ACCEPT'
uci set firewall.@zone[0].forward='ACCEPT'
uci set firewall.@zone[0].network='lan'
uci set firewall.@zone[1]=zone
uci set firewall.@zone[1].name='wan'
uci set firewall.@zone[1].input='REJECT'
uci set firewall.@zone[1].output='ACCEPT'
uci set firewall.@zone[1].forward='REJECT'
uci set firewall.@zone[1].masq='1'
uci set firewall.@zone[1].mtu_fix='1'
uci set firewall.@zone[1].network='wan wan6'
uci set firewall.@zone[3]=zone
uci set firewall.@zone[3].forward='REJECT'
uci set firewall.@zone[3].name='mesh'
uci set firewall.@zone[3].input='REJECT'
uci set firewall.@zone[3].output='REJECT'
uci set firewall.@zone[3].network='mesh_cable mesh_wifi'
uci set firewall.@zone[1]=zone
uci set firewall.@zone[1].name='wan'
uci set firewall.@zone[1].input='REJECT'
uci set firewall.@zone[1].output='ACCEPT'
uci set firewall.@zone[1].forward='REJECT'
uci set firewall.@zone[1].masq='1'
uci set firewall.@zone[1].mtu_fix='1'
uci set firewall.@zone[1].network='wan wan6'
uci set network.guest=interface
uci set network.guest.type='bridge'
uci set network.guest._orig_ifname='appub24'
uci set network.guest._orig_bridge='true'
uci set network.guest.ifname='bat0 appub24'
uci set network.guest.proto='dhcpv6'
uci set network.guest.reqaddress='try'
uci set network.guest.reqprefix='auto'
uci set network.mesh_cable=interface
uci set network.mesh_cable.ifname='eth1.11'
uci set network.mesh_cable.mtu='1528'
uci set network.mesh_cable.proto='batadv'
uci set network.mesh_cable.mesh='bat0'
uci set network.mesh_wifi=interface
uci set network.mesh_wifi.ifname='adhoc24'
uci set network.mesh_wifi.mtu='1528'
uci set network.mesh_wifi.proto='batadv'
uci set network.mesh_wifi.mesh='bat0'
uci set network.@switch[0]=switch
uci set network.@switch[0].name='switch0'
uci set network.@switch[0].reset='1'
uci set network.@switch[0].enable_vlan='1'
uci set network.@switch[0].mirror_source_port='0'
uci set network.@switch[0].mirror_monitor_port='0'
uci set network.@switch_vlan[2]=switch_vlan
uci set network.@switch_vlan[2].device='switch0'
uci set network.@switch_vlan[2].vlan='3'
uci set network.@switch_vlan[2].vid='11'
uci set network.@switch_vlan[2].ports='0t 5t'
uci set network.@switch_vlan[3]=switch_vlan
uci set network.@switch_vlan[3].device='switch0'
uci set network.@switch_vlan[3].vlan='4'
uci set network.@switch_vlan[3].vid='12'
uci set network.@switch_vlan[3].ports='0t 5t'
uci set network.@switch_vlan[0]=switch_vlan
uci set network.@switch_vlan[0].device='switch0'
uci set network.@switch_vlan[0].vlan='1'
uci set network.@switch_vlan[0].vid='3'
uci set network.@switch_vlan[0].ports='1 5t 6'
uci set network.@switch_vlan[1]=switch_vlan
uci set network.@switch_vlan[1].device='switch0'
uci set network.@switch_vlan[1].vlan='2'
uci set network.@switch_vlan[1].vid='5'
uci set network.@switch_vlan[1].ports='0t 2 3 4 5t'
uci set network.wan6=interface
uci set network.wan6.ifname='eth0'
uci set network.wan6.proto='dhcpv6'
uci set network.wan=interface
uci set network.wan.proto='dhcp'
uci set network.wan._orig_ifname='eth0'
uci set network.wan._orig_bridge='true'
uci set network.wan.ifname='eth0'
uci set dhcp.odhcpd=odhcpd
uci set dhcp.odhcpd.maindhcp='0'
uci set dhcp.odhcpd.leasefile='/tmp/hosts/odhcpd'
uci set dhcp.odhcpd.leasetrigger='/usr/sbin/odhcpd-update'
uci set system.ntp=timeserver
uci set system.ntp.server='0.openwrt.pool.ntp.org' '1.openwrt.pool.ntp.org' '2.openwrt.pool.ntp.org' '3.openwrt.pool.ntp.org'
uci set system.ntp.enabled='1'
uci set system.ntp.enable_server='1'
uci set system.@system[0].timezone='CET-1CEST,M3.5.0,M10.5.0/3'
uci set uhttpd.main=uhttpd
uci set uhttpd.main.listen_http='0.0.0.0:80' '[::]:80'
uci set uhttpd.main.listen_https='0.0.0.0:443' '[::]:443'
uci set uhttpd.main.redirect_https='1'
uci set uhttpd.main.home='/www'
uci set uhttpd.main.rfc1918_filter='1'
uci set uhttpd.main.max_requests='3'
uci set uhttpd.main.max_connections='100'
uci set uhttpd.main.cert='/etc/uhttpd.crt'
uci set uhttpd.main.key='/etc/uhttpd.key'
uci set uhttpd.main.cgi_prefix='/cgi-bin'
uci set uhttpd.main.script_timeout='60'
uci set uhttpd.main.network_timeout='30'
uci set uhttpd.main.http_keepalive='20'
uci set uhttpd.main.tcp_keepalive='1'
uci set uhttpd.main.ubus_prefix='/ubus'
uci set uhttpd.px5g=cert
uci set uhttpd.px5g.days='730'
uci set uhttpd.px5g.bits='1024'
uci set uhttpd.px5g.country='DE'
uci set uhttpd.px5g.state='Berlin'
uci set uhttpd.px5g.location='Berlin'
uci set uhttpd.px5g.commonname='OpenWrt'
uci set vnstat.@vnstat[0]=vnstat
uci set vnstat.@vnstat[0].interface='adhoc24' 'appriv24' 'appriv5' 'appub24' 'bat0' 'br-guest' 'br-lan'
uci set wireless.adhoc24=wifi-iface
uci set wireless.adhoc24.device='radio1'
uci set wireless.adhoc24.network='mesh'
uci set wireless.adhoc24.encryption='none'
uci set wireless.adhoc24.ssid='wifimesh-g-ac'
uci set wireless.adhoc24.mode='adhoc'
uci set wireless.adhoc24.ifname='adhoc24'
uci set wireless.adhoc24.bssid='CA:FF:24:FF:AC:24'
uci set wireless.adhoc24.disabled='0'
uci set wireless.adhoc24.mcast_rate='12000'
uci set wireless.appub24=wifi-iface
uci set wireless.appub24.device='radio1'
uci set wireless.appub24.mode='ap'
uci set wireless.appub24.ifname='appub24'
uci set wireless.appub24.ssid='Freifunk'
uci set wireless.appub24.encryption='none'
uci set wireless.appub24.network='guest'
uci set wireless.radio0=wifi-device
uci set wireless.radio0.type='mac80211'
uci set wireless.radio0.hwmode='11a'
uci set wireless.radio0.path='pci0000:01/0000:01:00.0'
uci set wireless.radio0.txpower='17'
uci set wireless.radio0.country='DE'
uci set wireless.radio0.htmode='VHT40'
uci set wireless.radio0.channel='44'
uci set wireless.radio1=wifi-device
uci set wireless.radio1.type='mac80211'
uci set wireless.radio1.hwmode='11g'
uci set wireless.radio1.path='platform/qca955x_wmac'
uci set wireless.radio1.country='DE'
uci set wireless.radio1.channel='11'
uci set wireless.radio1.txpower='20'
uci set network.lan=interface
uci set network.lan.type='bridge'
uci set network.lan._orig_ifname='eth1 appriv'
uci set network.lan._orig_bridge='true'
uci set network.lan.proto='static'
uci set network.lan.ifname='appriv24 appriv5 eth1.5'
uci set network.lan.netmask='255.255.255.0'
uci set network.lan.ipaddr='10.10.101.46'
uci set network.lan.gateway='10.10.101.49'
uci set wireless.appriv24=wifi-iface
uci set wireless.appriv24.device='radio1'
uci set wireless.appriv24.mode='ap'
uci set wireless.appriv24.ifname='appriv24'
uci set wireless.appriv24.network='lan'
uci set wireless.appriv24.ssid='Spaltung-legacy'
uci set wireless.appriv24.encryption='wpa2+ccmp'
uci set wireless.appriv24.auth_server='10.10.101.49'
uci set wireless.appriv24.acct_server='10.10.101.49'
uci set wireless.appriv24.auth_secret='mRePVJc5m4GsrE4OEJKL'
uci set wireless.appriv24.acct_secret='mRePVJc5m4GsrE4OEJKL'
uci set wireless.appriv24.nasid='c7groundfloor'
uci set wireless.appriv5=wifi-iface
uci set wireless.appriv5.device='radio0'
uci set wireless.appriv5.network='lan'
uci set wireless.appriv5.mode='ap'
uci set wireless.appriv5.ifname='appriv5'
uci set wireless.appriv5.ssid='Spaltung'
uci set wireless.appriv5.encryption='wpa2+ccmp'
uci set wireless.appriv5.auth_server='10.10.101.49'
uci set wireless.appriv5.acct_server='10.10.101.49'
uci set wireless.appriv5.nasid='c7groundfloor'
uci set wireless.appriv5.auth_secret='mRePVJc5m4GsrE4OEJKL'
uci set wireless.appriv5.acct_secret='mRePVJc5m4GsrE4OEJKL'
uci commit
exit 0
