# files/etc/uci-defaults/99_custom_settings
#!/bin/ash
#
# env variables
hostnamebasis="ffcp-"
randomized=$(</dev/urandom tr -dc A-Z-a-z-0-9 | head -c24)
sourceeth="eth0"
ulafirst8="FD"
guestnetidentifier=1
privnetidentifier=2

#get last 40 bits from $sourceeth interface (spine)
hwaddr=$(ifconfig $sourceeth|grep HWaddr|awk '{print $5}')
hwaddrwocolon=$(echo $hwaddr|sed -e 's/://g')
guestnet=`echo $ulafirst8${hwaddrwocolon:2:2}:${hwaddrwocolon:4:2}${hwaddrwocolon:6:2}:${hwaddrwocolon:8:2}${hwaddrwocolon:10:1}$guestnetidentifier::/64)`
guestnetip=`echo $ulafirst8${hwaddrwocolon:2:2}:${hwaddrwocolon:4:2}${hwaddrwocolon:6:2}:${hwaddrwocolon:8:2}${hwaddrwocolon:10:1}$guestnetidentifier::1/64)`
privnet=`echo $ulafirst8${hwaddrwocolon:2:2}:${hwaddrwocolon:4:2}${hwaddrwocolon:6:2}:${hwaddrwocolon:8:2}${hwaddrwocolon:10:1}$privnetidentifier::/64)`
privnetip=`echo $ulafirst8${hwaddrwocolon:2:2}:${hwaddrwocolon:4:2}${hwaddrwocolon:6:2}:${hwaddrwocolon:8:2}${hwaddrwocolon:10:1}$privnetidentifier::1/64)`

# UCI
uci set system.@system[0].hostname=$(echo -n "ffcp-$hwaddrwocolon")
uci set system.ntp.enable_server=1
uci set dhcp.@dnsmasq[0].domain=$(echo -n "ffcp.lan")
uci set system.@system[0].timezone="CET-1CEST,M3.5.0,M10.5.0/3"

# set release/version
grep Version /etc/banner | awk '{ print $2}' > /lib/gluon/gluon-version
grep Release /etc/banner | awk '{ print $2}' > /lib/gluon/release

# set default MAC for gluon
echo $hwaddr > /lib/gluon/core/sysconfig/primary_mac

# repair init script announced
#sed -i 's/collect.lua/announce.lua/g' /etc/hotplug.d/iface/10-gluon-announced
chmod +x /lib/gluon/announce/collect.lua

exit 0
