# files/etc/uci-defaults/99_custom_settings
../../../../filesnkac/etc/uci-defaults/99_custom_settings